# dist/shortly.service
[Unit]
Description=shortly
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/shortly
ExecStart=/srv/shortly/app
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=Shortly

[Install]
WantedBy=multi-user.target
