# dist/shortly.service
[Unit]
Description=shortly
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/shortly
# ExecStart=/srv/shortly/bin/app  >>/srv/shortly/app.log  2>&1
ExecStart=/srv/shortly/bin/app
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=Shortly

[Install]
WantedBy=multi-user.target
